# weather.timer
[Unit]
Description=A timer to run the weather script

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
